# linux-rt-with-192kHz-souind-driver/ti-processor-sdk-linux-rt-am62dxx-evm-11.02.08.02/linux-devkit/sysroots/aarch64-oe-linux/var/lib/opkg/info/avahi-daemon.list
/usr/lib/systemd/system-preset/98-avahi-daemon.preset	0100644
/usr/lib/systemd/system/avahi-daemon.socket	0100644
/usr/share/dbus-1/interfaces	040755
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.DomainBrowser.xml	0100644
/usr/lib/systemd/system	040755
/usr/share/dbus-1/system-services	040755
/usr/lib	040755
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.ServiceBrowser.xml	0100644
/usr/share/dbus-1	040755
/usr/sbin/avahi-daemon	0100755
/usr/share/avahi	040755
/usr/sbin	040755
/etc/dbus-1/system.d/avahi-dbus.conf	0100644
/usr/share/avahi/avahi-service.dtd	0100644
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.HostNameResolver.xml	0100644
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.ServiceResolver.xml	0100644
/usr/share/dbus-1/system-services/org.freedesktop.Avahi.service	0100644
/etc/avahi/avahi-daemon.conf	0100644
/usr	040755
/etc	040755
/etc/dbus-1/system.d	040755
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.EntryGroup.xml	0100644
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.RecordBrowser.xml	0100644
/etc/avahi/hosts	0100644
/usr/lib/systemd	040755
/usr/lib/systemd/system/avahi-daemon.service	0100644
/etc/avahi/services	040755
/usr/lib/systemd/system-preset	040755
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.Server.xml	0100644
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.ServiceTypeBrowser.xml	0100644
/usr/share	040755
/etc/dbus-1	040755
/usr/share/dbus-1/interfaces/org.freedesktop.Avahi.AddressResolver.xml	0100644
/etc/avahi	040755
