# linux-rt-with-192kHz-souind-driver/ti-processor-sdk-linux-rt-am62dxx-evm-11.02.08.02/linux-devkit/sysroots/aarch64-oe-linux/var/lib/opkg/info/bluez5.list
/usr/lib/systemd/system/bluetooth.service	0100644
/usr/share/dbus-1/system.d/bluetooth.conf	0100644
/usr/lib/systemd/system	040755
/etc/bluetooth/network.conf	0100644
/etc/bluetooth/input.conf	0100644
/usr/share/dbus-1/system-services	040755
/usr/bin/mpris-proxy	0100755
/usr/lib	040755
/etc/bluetooth	040755
/usr/libexec	040755
/usr/share/dbus-1	040755
/usr/bin/btattach	0100755
/usr/share/dbus-1/system-services/org.bluez.service	0100644
/usr/bin/hciconfig	0100755
/usr/bin/hcidump	0100755
/usr/lib/libbluetooth.so.3.19.12	0100755
/usr/bin/ciptool	0100755
/usr/lib/systemd/system-preset/98-bluez5.preset	0100644
/usr/bin/l2ping	0100755
/usr/bin/isotest	0100755
/usr/bin/hcitool	0100755
/usr/bin	040755
/usr/libexec/bluetooth	040755
/usr/bin/bluemoon	0100755
/usr/bin/sdptool	0100755
/usr/bin/hex2hcd	0100755
/usr/bin/bluetoothctl	0100755
/usr	040755
/usr/bin/hciattach	0100755
/etc	040755
/usr/libexec/bluetooth/bluetoothd	0100755
/usr/lib/systemd	040755
/usr/bin/l2test	0100755
/usr/lib/systemd/system-preset	040755
/usr/lib/libbluetooth.so.3	0120777	libbluetooth.so.3.19.12
/usr/bin/rctest	0100755
/usr/share	040755
/usr/share/dbus-1/system.d	040755
/usr/bin/btmon	0100755
/usr/bin/rfcomm	0100755
