# linux-rt-with-192kHz-souind-driver/ti-processor-sdk-linux-rt-am62dxx-evm-11.02.08.02/linux-devkit/sysroots/aarch64-oe-linux/var/lib/opkg/info/libpam-runtime.list
/usr/lib/systemd/system/pam_namespace.service	0100644
/etc/security/faillock.conf	0100644
/etc/security/group.conf	0100644
/usr/sbin/unix_chkpwd	0104755
/usr/lib/systemd/system	040755
/usr/lib	040755
/usr/sbin/unix_update	0100755
/etc/pam.d/common-password	0100644
/usr/sbin	040755
/etc/pam.d/common-session	0100644
/etc/security/namespace.init	0100755
/etc/security/limits.conf	0100644
/usr/sbin/mkhomedir_helper	0100755
/usr/sbin/pam_timestamp_check	0100755
/usr/sbin/pwhistory_helper	0100755
/usr	040755
/etc	040755
/etc/pam.d/common-auth	0100644
/usr/lib/systemd	040755
/etc/security/access.conf	0100644
/etc/security/pwhistory.conf	0100644
/etc/security/time.conf	0100644
/etc/pam.d/common-session-noninteractive	0100644
/etc/pam.d	040755
/etc/security/pam_env.conf	0100644
/etc/security	040755
/etc/security/limits.d	040755
/etc/pam.d/other	0100644
/usr/sbin/pam_namespace_helper	0100755
/etc/security/namespace.d	040755
/etc/pam.d/common-account	0100644
/etc/security/namespace.conf	0100644
/usr/sbin/faillock	0100755
