# linux-rt-with-192kHz-souind-driver/ti-processor-sdk-linux-rt-am62dxx-evm-11.02.08.02/linux-devkit/sysroots/aarch64-oe-linux/var/lib/opkg/info/mesa-megadriver.list
/usr/lib/dri/virtio_gpu_dri.so	0100755
/usr/lib/dri	040755
/usr/share/drirc.d/00-mesa-defaults.conf	0100644
/usr/lib/dri/r300_dri.so	0100755
/usr/lib/dri/swrast_dri.so	0100755
/usr/lib/dri/kms_swrast_dri.so	0100755
/usr/lib/dri/zink_dri.so	0100755
/usr/lib/dri/nouveau_dri.so	0100755
/usr/share/drirc.d	040755
