# linux-rt-with-192kHz-souind-driver/ti-processor-sdk-linux-rt-am62dxx-evm-11.02.08.02/linux-devkit/sysroots/aarch64-oe-linux/var/lib/opkg/info/ofono.list
/usr/lib/systemd/system	040755
/usr/lib	040755
/usr/sbin	040755
/etc/ofono	040755
/usr/share/ofono/provision.db	0100644
/etc/dbus-1/system.d/ofono.conf	0100644
/usr	040755
/etc	040755
/usr/sbin/ofonod	0100755
/usr/share/ofono	040755
/etc/dbus-1/system.d	040755
/usr/lib/systemd	040755
/etc/ofono/phonesim.conf	0100644
/usr/lib/systemd/system-preset	040755
/usr/share	040755
/etc/dbus-1	040755
/usr/lib/systemd/system/ofono.service	0100644
/usr/lib/systemd/system-preset/98-ofono.preset	0100644
