# linux-rt-with-192kHz-souind-driver/ti-processor-sdk-linux-rt-am62dxx-evm-11.02.08.02/linux-devkit/sysroots/aarch64-oe-linux/var/lib/opkg/info/systemd-extra-utils.list
/usr/bin/systemd-cat	0100755
/usr/bin/systemd-stdio-bridge	0100755
/usr/lib/systemd/system/systemd-initctl.service	0100644
/usr/lib/systemd/system	040755
/usr/bin/systemd-delta	0100755
/usr/sbin/mount.ddi	0120777	../bin/systemd-dissect
/usr/lib	040755
/usr/lib/systemd/system-generators/systemd-gpt-auto-generator	0100755
/usr/lib/systemd/system/systemd-hybrid-sleep.service	0100644
/usr/bin/systemd-creds	0100755
/usr/sbin	040755
/usr/lib/systemd/system/systemd-initctl.socket	0100644
/usr/bin/systemd-dissect	0100755
/usr/lib/systemd/systemd-socket-proxyd	0100755
/usr/lib/systemd/systemd-initctl	0100755
/usr/lib/systemd/systemd-sleep	0100755
/usr/bin	040755
/usr/bin/systemd-cgtop	0100755
/usr/lib/systemd/system/sockets.target.wants/systemd-initctl.socket	0120777	../systemd-initctl.socket
/usr/lib/systemd/system/sleep.target	0100644
/usr	040755
/usr/lib/systemd/system/systemd-suspend.service	0100644
/usr/bin/systemd-detect-virt	0100755
/usr/lib/systemd/system-sleep	040755
/usr/lib/systemd	040755
/usr/bin/systemd-run	0100755
/usr/bin/systemd-escape	0100755
/usr/bin/systemd-inhibit	0100755
/usr/lib/systemd/system-generators	040755
/usr/lib/systemd/systemd-cgroups-agent	0100755
/usr/lib/systemd/system/sockets.target.wants	040755
/usr/bin/systemd-cgls	0100755
/usr/bin/systemd-path	0100755
/usr/lib/systemd/system/systemd-hibernate.service	0100644
