# linux-rt-with-192kHz-souind-driver/ti-processor-sdk-linux-rt-am62dxx-evm-11.02.08.02/linux-devkit/sysroots/aarch64-oe-linux/var/lib/opkg/info/systemd-vconsole-setup.list
/usr/lib/systemd/system	040755
/usr/lib	040755
/usr/lib/systemd/system/systemd-vconsole-setup.service	0100644
/usr/lib/systemd/systemd-vconsole-setup	0100755
/usr	040755
/usr/lib/systemd	040755
